# server/etc/cron.daily/steeplechase-long
#!/bin/sh

[ ! -e ~/logs ] && mkdir ~/logs

`cd ~/logs && logrotate -f -s ~/etc/logrotate.state ~/etc/steeplechase-long.rotate.conf` 
~/bin/steeplechase-long > ~/logs/steeplechase-long.log 2>&1
