# server/etc/cron.daily/update-nightly
#!/bin/sh

[ ! -e ~/logs ] && mkdir ~/logs

`cd ~/logs && logrotate -f -s ~/etc/logrotate.state ~/etc/update-nightly.rotate.conf`
~/bin/update-nightly > ~/logs/update-nightly.log 2>&1
